# packaging/nava.service.in
[Unit]
Description=NAVA voice dictation daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=@NAVA@ _run
Restart=on-failure
RestartSec=3
# X11 is required for injection + hotkey; templated by install.sh from the install session.
Environment=DISPLAY=@DISPLAY@

[Install]
WantedBy=default.target
